# one-api.service
# File path: /etc/systemd/system/one-api.service
# sudo systemctl daemon-reload
# sudo systemctl start one-api
# sudo systemctl enable one-api
# sudo systemctl status one-api
[Unit]
Description=One API Service
After=network.target

[Service]
User=root
WorkingDirectory=/www/wwwroot/api.f91.dev
ExecStart=/www/wwwroot/api.f91.dev/one-api --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
